# roles/backups/templates/docker-backup.service.j2
[Unit]
Description=Backup Docker configuration and manifests
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_docker.sh
