# roles/backups/templates/docker-backup.timer.j2
[Unit]
Description=Run Docker backup script

[Timer]
OnCalendar={{ backup_timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
